# Antd2/Targets/antd.target
[Unit]
Description=Anthilla OS - antd Target

[Install]
WantedBy=basic.target