# Antd2/Targets/applicative.target
[Unit]
Description=Anthilla OS - Applicative Target

[Install]
WantedBy=basic.target
